# templates/etc/init/autossh-tunnel-client.conf.j2
# {{ ansible_managed }}

# autossh-tunnel-client

description "Set up a persistent tunnel (using autossh)"
author      "Mischa ter Smitten"

start on net-device-up IFACE=eth0
stop on runlevel [01S6]

respawn
respawn limit 5 60 # Respawn max 5 times in 60 seconds
umask 022

pre-start script
  test -x /usr/bin/autossh || { stop; exit 0; }

  mkdir -p -m0755 {{ autossh_tunnel_client_autossh_pidfile | dirname }}
end script

script
  export AUTOSSH_DEBUG={{ autossh_tunnel_client_autossh_debug }}
  export AUTOSSH_FIRST_POLL={{ autossh_tunnel_client_autossh_first_poll }}
  export AUTOSSH_GATETIME={{ autossh_tunnel_client_autossh_gatetime }}
  export AUTOSSH_LOGLEVEL={{ autossh_tunnel_client_autossh_loglevel }}
  export AUTOSSH_PIDFILE={{ autossh_tunnel_client_autossh_pidfile }}
  export AUTOSSH_POLL={{ autossh_tunnel_client_autossh_poll }}

  autossh -M 0 -4 -N -L {{ autossh_tunnel_client_forward }} {{ autossh_tunnel_client_user }}@{{ autossh_tunnel_client_host }} \
    -p {{ autossh_tunnel_client_port }} \
    -i {{ autossh_tunnel_client_configuration_directory }}/{{ autossh_tunnel_client_identity }} \
{% for ssh_options in autossh_tunnel_client_ssh_options %}
    -o "{{ ssh_options }}" \
{% endfor %}
  ;
end script
